# supervisord/supervisord.conf
[supervisord]
nodaemon=true

[supervisorctl]

[program:php-fpm]
command=rc-service php-fpm7 restart

[program:nginx]
command=rc-service nginx restart